# MoarGif.service
[Unit]
Description=MoarGif

[Service]
User=pi
ExecStart=/home/pi/App/node-v4.6.0-linux-armv6l/bin/node /home/pi/App/MoarGif/server.js config=rpi-config.json
Restart=on-failure
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
